# etc/conf.d/celery.conf
[program:celery]
command = ./manage.py celery worker --loglevel=info
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
