# etc/conf.d/django.conf
[program:django]
command = gunicorn --bind 0.0.0.0:8000 skeleton.wsgi
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
